# unit-files/generate-index.timer
[Unit]
Description=Run generate-index.service every day at 05:00:00
[Timer]
OnCalendar=*-*-* 05:00:00
[Install]
WantedBy=timers.target
